# claw.service
[Unit]
After=network.target

[Service]
ExecStart=claw

[Install]
WantedBy=multi-user.target
