# osprey.service
[Unit]
After=network.target

[Service]
ExecStart=osprey

[Install]
WantedBy=multi-user.target
